# templates/deb/buildbox-agent.upstart
description "buildbox-agent"

respawn

start on runlevel [2345]
stop on runlevel [06]

script
source /etc/buildbox-agent/buildbox-agent.env
mkdir -p $BUILDBOX_DIR
chown -R $BUILDBOX_USER:$BUILDBOX_USER_GROUP $BUILDBOX_DIR
su - $BUILDBOX_USER -c "buildbox-agent start"
end script
