# templates/deb/buildbox-agent.upstart
description "buildbox-agent"

respawn

start on runlevel [2345]
stop on runlevel [06]

script
su - root -c "buildbox-agent start"
end script
