# templates/deb/buildbox-agent.upstart
description "buildbox-agent"

respawn

start on runlevel [2345]
stop on runlevel [06]

kill signal TERM
reload signal USR2

env PIDFILE=/var/run/buildbox-agent.pid

pre-start script
. /etc/buildbox-agent/buildbox-agent.env
mkdir -p $BUILDBOX_DIR
chown -R $BUILDBOX_USER:$BUILDBOX_USER_GROUP $BUILDBOX_DIR
end script

script
set -a
. /etc/buildbox-agent/buildbox-agent.env
set +a
exec start-stop-daemon --start --quiet --chuid $BUILDBOX_USER --pidfile $PIDFILE --make-pidfile --exec /usr/bin/buildbox-agent -- start
end script

post-stop exec rm -f $PIDFILE
