# templates/buildbox.conf
respawn
respawn limit 15 5

start on runlevel [2345]
stop on runlevel [06]

script
su - ubuntu -c "buildbox-agent start --access-token token1234"
end script
